# crontab
*/5 * * * * root "cd app && ./locomotive -vvv $(printenv REMOTE_SERVER) >> /var/log/cron.log 2>&1"